# roles/tor/templates/polipo-supervisord.conf.j2
[program:polipo]
command=/usr/bin/polipo -c /etc/polipo/config
autostart=true
autorestart=true
startsecs=10
startretries=3