# roles/tor/templates/tor-supervisord.conf.j2
[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc
autostart=true
autorestart=true
startsecs=10
startretries=3